# robot-console-backend.service
[Unit]
Description=Robot Console Backend API
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=exec
User=robot-console
Group=robot-console
WorkingDirectory=/opt/robot-console/backend
Environment=PATH=/opt/robot-console/backend/venv/bin
ExecStart=/opt/robot-console/backend/venv/bin/gunicorn main:app -w 4 -k uvicorn.workers.UvicornWorker --bind 0.0.0.0:8000
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/robot-console/backend/logs /opt/robot-console/backend/videos

# Environment file
EnvironmentFile=/opt/robot-console/backend/.env

[Install]
WantedBy=multi-user.target